# projects/minecraft/ops/friedutchplus.minecraft.server.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>friedutchplus.minecraft.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/administrator/Sites/friedutchplus/projects/minecraft/ops/start_server.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SERVER_ROOT</key>
    <string>/Users/administrator/Servers/minecraft</string>
    <key>JAVA_BIN</key>
    <string>/opt/homebrew/opt/openjdk/bin/java</string>
    <key>JAR_NAME</key>
    <string>paper.jar</string>
    <key>MIN_RAM</key>
    <string>2G</string>
    <key>MAX_RAM</key>
    <string>4G</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/administrator/Servers/minecraft</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/administrator/Library/Logs/friedutchplus.minecraft.server.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/administrator/Library/Logs/friedutchplus.minecraft.server.err.log</string>
</dict>
</plist>
